# god-orchestrator.service
[Unit]
Description=GodStack Production Ingestion Matrix Daemon
After=network.target

[Service]
Type=simple
User=tangleroot013
WorkingDirectory=/home/tangleroot013/god_stack
ExecStart=/usr/bin/python3 run_production_matrix.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/tangleroot013/god_stack/logs/daemon_stdout.log
StandardError=append:/home/tangleroot013/god_stack/logs/daemon_stderr.log

[Install]
WantedBy=multi-user.target
