# god-orchestrator.service
[Unit]
Description=GodStack Production Ingestion Matrix Daemon
After=network.target

[Service]
Type=simple
User=tangleroot013
Group=tangleroot013
WorkingDirectory=/home/tangleroot013/god_stack

# Production Environment Adjustments
Environment=GOD_STACK_ENV=production
Environment=PYTHONUNBUFFERED=1
Environment=PROMETHEUS_PORT=8015

# Review Hardening: Protect process visibility and control process capabilities limits
ProtectSystem=full
ProtectHome=yes
ProtectProc=invisible
AmbientCapabilities=
LimitNOFILE=4096
NoNewPrivileges=true
PrivateTmp=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

ExecStart=/home/tangleroot013/god_stack/.venv/bin/python3 run_production_matrix.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/tangleroot013/god_stack/logs/daemon_stdout.log
StandardError=append:/home/tangleroot013/god_stack/logs/daemon_stderr.log

[Install]
WantedBy=multi-user.target
